# ROS/mapping/mapping.service
[Unit]
Description= Lancement de la cartographie. Mapping's launch

[Service]
ExecStart=/<chemin/pathway>/mapping.sh
