# ROS/mapping/rosmaster.service
[Unit]
Description= lancement du rosmaster. Rosmaster's launch

[Service]
ExecStart=/<chemin/pathway>/rosmaster.sh
