# ROS/mapping/rosmaster.service
[Unit]
Description= lancement du rosmaster. Rosmaster's launch

[Service]
ExecStart=/etc/systemd/system/rosmaster.sh
