# ROS/mapping/saveamap.service
[Unit]
Description= lancement de la sauvegarde. Save's launch

[Service]
ExecStart=/<chemin/pathway>/saveamap.sh
