# Vermot-Desroches Matthias/NVIDIA/TD-9-11-22/dump-database.service
[Unit]
Description= backup database

[Service]
ExecStart=/<chemin>/dump-database.sh
