# Vermot-Desroches Matthias/NVIDIA/TD-9-11-22/dump-database.timer
[Unit]
Description= run backup every 10 mins

[Timer]
OnCalendar=*:0/10
AccuracySec=1s
Unit=dump-database.service

[Install]
WantedBy=timers.target
